# conf/rc.local
#!/bin/sh -e
cd /var/www/facetrack
leds/blinkip &
nodejs capture.js &
exit 0
