# conf/rc.local
#!/bin/sh -e
cd /var/www/robot
leds/blinkip &
nodejs facetrack.js &
exit 0
